# profile.d/dendyops_other.sh
HISTSIZE=300000
HISTTIMEFORMAT="[ %F %T ] "
export HISTSIZE
export HISTTIMEFORMAT
export PATH=$PATH:$HOME/bin:/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin:/opt/nginx/sbin/
